# deploy/berries-ingest.service
[Unit]
Description=Berries Ingest API
After=network.target

[Service]
User=berries
WorkingDirectory=/opt/berries
ExecStart=/opt/berries/venv/bin/uvicorn ingest_api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/berries/.env

# ── Sandboxing / hardening ────────────────────────────────────────────────
# If any of these break the service on your host, remove the offending line
# and re-run `systemd-analyze security berries-ingest` to check overall score.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/berries/data /opt/berries/logs
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# NOTE: The previous deny list (`~@privileged @resources`) killed the service
# with SIGSYS because PyTorch/sentence-transformers (pulled in by ChromaDB)
# call syscalls in @resources during init. The @system-service allowlist alone
# is already a reasonable tightening; the deny list is omitted deliberately.
UMask=0077

[Install]
WantedBy=multi-user.target
