# deploy/berries-utils.service
[Unit]
Description=Berries Stream Utils
After=network.target berries-ingest.service

[Service]
User=berries
WorkingDirectory=/opt/berries
ExecStart=/opt/berries/venv/bin/uvicorn stream_utils.main:app --host 127.0.0.1 --port 8002
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/berries/.env

[Install]
WantedBy=multi-user.target
